# share/docker/etc/services.d/11-nginx/run
#!/bin/sh -e

echo "[services.d] starting service: nginx"
nginx -g "daemon off;"
